# system_files/shared/usr/lib/systemd/user/reboot-notifier-show.service
[Unit]
Description=Show reboot notification to user

[Service]
Type=oneshot
ExecStart=/usr/libexec/ublue-reboot-notifier-show.sh
