# system_files/shared/usr/lib/systemd/system/taidan-initial-setup-reconfiguration.service.d/override.conf
[Unit]
ConditionPathExists=
ConditionPathExists=/var/lib/taidan
